# crontab
0 12 * * 5 cd /srv/sms-backup-service && yarn start
