# crontab
0 14 * * * node /srv/sms-backup-service/dist/index.js
