# crontab
0 12 * * 5 node /srv/sms-backup-service/dist/index.js
